# deploy/systemd/memoryx-api.service
[Unit]
Description=MemoryX API Server
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=exec
User=memoryx
Group=memoryx
WorkingDirectory=/data/memoryx/api

# Environment
Environment="PYTHONPATH=/data/memoryx/api"
Environment="DATABASE_URL=postgresql://memoryx:password@localhost:5432/memoryx"
Environment="REDIS_URL=redis://localhost:6379/0"
Environment="SECRET_KEY=your-secret-key-here"
Environment="OLLAMA_HOST=http://192.168.31.65:11434"

# Load additional env from file
EnvironmentFile=-/etc/memoryx/api.env

# Start command
ExecStart=/usr/local/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# Graceful shutdown
TimeoutStopSec=30
KillSignal=SIGTERM

# Logging
StandardOutput=append:/var/log/memoryx/api.log
StandardError=append:/var/log/memoryx/api-error.log

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/data/memoryx /var/log/memoryx

[Install]
WantedBy=multi-user.target
