# deploy/systemd/memoryx-webhook.service
[Unit]
Description=MemoryX Webhook Receiver
After=network.target

[Service]
Type=simple
User=memoryx
Group=memoryx
WorkingDirectory=/data/memoryx

# Environment
Environment="DEPLOY_TOKEN=your-secret-token-here"
Environment="PYTHONPATH=/data/memoryx/deploy/scripts"

# Start command
ExecStart=/usr/bin/python3 /data/memoryx/deploy/scripts/webhook_server.py

# Restart policy
Restart=always
RestartSec=5

# Logging
StandardOutput=append:/var/log/memoryx/webhook.log
StandardError=append:/var/log/memoryx/webhook-error.log

[Install]
WantedBy=multi-user.target
