# deploy/systemd/memoryx-api-docker.service
[Unit]
Description=MemoryX API Docker Container
Requires=docker.service
After=docker.service

[Service]
Type=simple
Restart=always
RestartSec=5

# 拉取最新镜像并运行
ExecStartPre=-/usr/bin/docker pull ghcr.io/t0ken-ai/memoryx-api:latest
ExecStartPre=-/usr/bin/docker rm -f memoryx-api

ExecStart=/usr/bin/docker run \
    --name memoryx-api \
    --restart=unless-stopped \
    -p 127.0.0.1:8000:8000 \
    -v /data/memoryx/static:/app/static:ro \
    -v /etc/memoryx/api.env:/app/.env:ro \
    --env-file /etc/memoryx/api.env \
    ghcr.io/t0ken-ai/memoryx-api:latest

ExecStop=/usr/bin/docker stop -t 30 memoryx-api
ExecStopPost=-/usr/bin/docker rm -f memoryx-api

[Install]
WantedBy=multi-user.target
